# base-image/workshop/rootfs/etc/systemd/system/setup-swap.service
[Unit]
Description=Setup swap

[Service]
Type=oneshot
ExecStart=/opt/misc/bin/setup-swap

[Install]
WantedBy=multi-user.target
