# base-image/workshop/rootfs/etc/profile.d/workshop.sh
source /opt/ros/humble/setup.bash

export ROS_DOMAIN_ID=21
